# overlay-squashfs-x86/etc/init.d/load-fonts-keymaps
#!/sbin/runscript

source /sbin/livecd-functions.sh

depend()
{
	after pwgen sshd portmap hald
	before sysresccd
}

start()
{
	ebegin "Loading fonts and keymaps"

	keymap=$(grep '^XKEYBOARD=' /etc/sysconfig/keyboard 2>/dev/null | cut -d= -f2)

	case "${keymap}" in
		fr|16)
			loadkeys fr-latin9
			;;
		uk|40)
			loadkeys uk
			;;
		ru|31)
			loadkeys ru
			setfont LatArCyrHeb-16
			;;
		pl|28)
			loadkeys pl
			setfont -m 8859-2
			;;
	esac
}

